# config/systemd/tradingbot-backup.timer
[Unit]
Description=Trading Bot Daily Backup Timer
Documentation=https://github.com/your-org/tradingbot
Requires=tradingbot-backup.service

[Timer]
# Run daily at 2:00 AM
OnCalendar=daily
OnCalendar=*-*-* 02:00:00

# Run 10 minutes after boot if missed
Persistent=true
OnBootSec=10min

# Randomize start time by up to 5 minutes to avoid load spikes
RandomizedDelaySec=5min

# Accuracy
AccuracySec=1min

[Install]
WantedBy=timers.target
